# autostart/ColourDepth32.desktop
[Desktop Entry]
Name=ColourDepth32
Exec="startx  -- -depth 32"
Type=Application
Terminal=true
Hidden=false

